# systemd/.config/systemd/user/backup.service
[Unit]
Description=Home directory backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-home
